# supervisord.conf
[unix_http_server]
file=/run/supervisord.sock  ; the path to the socket file

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///run/supervisord.sock ; use a unix:// URL for a unix socket

[supervisord]
nodaemon=true

[program:redis]
directory=/
user=root
command=redis-server
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:postgres]
directory=/tmp/db/postgres
user=postgres
command=postgres -D /tmp/db/postgres
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:init]
directory=/
user=root
command=/init.sh
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
autorestart=false

[program:stack]
directory=/
user=thethings
command=ttn-lw-stack -c /config/ttn-lw-stack-docker.yml start
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
autostart=false
